# Linux/s1.service
[Unit]
Description=CPU Stress Command Testing

[Service]
ExecStart= dd if=/dev/zero of=/dev/null
MemoryLimit=100M
CPUShares=1024