# Linux/s2.service
[Unit]
Description=CPU Stress Command Testing2

[Service]
ExecStart= dd if=/dev/zero of=/dev/null
MemoryLimit=100M
CPUShares=2048